# root/etc/services.d/ssh-mitm/run
#!/command/with-contenv bash

exec ssh-proxy --lport=4444 --rhost=localhost --rport=22 localhost
